# docker/supervisor/queue-worker.conf
[program:menu-queue-worker]
command=php /var/www/html/artisan queue:work redis --tries=3 --timeout=120 --sleep=3 --max-jobs=500 --max-time=3600
directory=/var/www/html
autostart=true
autorestart=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/log/supervisor/menu-worker.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stopwaitsecs=130
